# deploy/kalshi-bot.service
[Unit]
Description=Kalshi AI Trading Bot
Documentation=https://github.com/Ron-Web-Dotcom/Kalshi-Trading-Bot
After=network-online.target
Wants=network-online.target
# Restart if the process dies for any reason
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=kalshi
Group=kalshi
WorkingDirectory=/opt/kalshi-bot
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/kalshi-bot/.env
ExecStart=/opt/kalshi-bot/.venv/bin/python /opt/kalshi-bot/bot.py

# Restart policy — always restart unless manually stopped
Restart=always
RestartSec=15

# Resource limits
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/kalshi-bot

# Logging — all output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kalshi-bot

[Install]
WantedBy=multi-user.target
